# systemd/renogy-bt.service
[Unit]
Description=Renogy BT poller and Home Assistant Bluetooth proxy
After=network-online.target bluetooth.target
Wants=network-online.target bluetooth.target

[Service]
Type=simple
User=jorik
Group=jorik
WorkingDirectory=/home/jorik/renogy-bt
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/jorik/renogy-bt/venv/bin/python /home/jorik/renogy-bt/ha_proxy_example.py /home/jorik/renogy-bt/config.ini
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target
